# shared_files/postgres.cron
SHELL=/bin/bash
MAILTO=""

0 0 * * * /backup/backup_philookup.sh > /backup/backup.log 2>&1

